# products/fortress/systemd/fortress-network.service
#!/usr/bin/env systemd
#
# fortress-network.service - Fortress Network Configuration (Filter Mode)
#
# Restores network configuration after reboot for filter mode:
# - VLAN 200 management interface (IP address)
# - WiFi port VLAN tags (for proper traffic flow)
# - LAN port VLAN configuration
#
# This service is for FILTER MODE (default). For VLAN mode, use fortress-vlan.service
#
# Install:
#   sudo cp fortress-network.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable fortress-network.service
#

[Unit]
Description=HookProbe Fortress Network Configuration (Filter Mode)
Documentation=https://hookprobe.com/docs/fortress
After=network-online.target openvswitch-switch.service
Wants=network-online.target
Requires=openvswitch-switch.service

# Run after hostapd so WiFi interfaces exist on bridge
After=fortress-hostapd-24ghz.service fortress-hostapd-5ghz.service
Wants=fortress-hostapd-24ghz.service fortress-hostapd-5ghz.service

[Service]
Type=oneshot
RemainAfterExit=yes

# Wait for OVS to be fully ready
ExecStartPre=/bin/sleep 3

# Run filter mode network setup
ExecStart=/opt/hookprobe/products/fortress/devices/common/filter-mode-network.sh setup

[Install]
WantedBy=multi-user.target
