# products/fortress/systemd/fts-fingerprint-engine.service
[Unit]
Description=HookProbe Fortress Unified Fingerprint Engine
Documentation=https://github.com/hookprobe/hookprobe
After=network.target openvswitch-switch.service fortress.service dnsmasq.service
# Only run if the Python module exists
ConditionPathExists=/opt/hookprobe/fortress/lib/unified_fingerprint_engine.py
# Start rate limiting (must be in [Unit] section)
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/hookprobe/fortress/lib/unified_fingerprint_engine.py --daemon
WorkingDirectory=/opt/hookprobe/fortress

# Run as root - needs access to OVS and dnsmasq lease file
# Security is maintained via ProtectSystem and other hardening
User=root
Group=root

# Environment - include lib directory in Python path for imports
Environment=PYTHONPATH=/opt/hookprobe/fortress/lib:/opt/hookprobe/fortress:/opt/hookprobe
Environment=FORTRESS_CONFIG=/etc/hookprobe/fortress.conf

# Restart policy
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fts-fingerprint-engine

# Security hardening (still applies even when running as root)
ProtectHome=true
PrivateTmp=true
# Note: NoNewPrivileges and ProtectSystem=strict removed as they conflict with OVS access

[Install]
WantedBy=multi-user.target
