# products/fortress/systemd/fts-lte-collector.service
[Unit]
Description=Fortress LTE Usage Tracker (Watermark-based)
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/hookprobe/fortress/devices/common/lte-usage-tracker.sh collect
User=root
Group=root

# Minimal resources
Nice=19
IOSchedulingClass=idle

[Install]
WantedBy=multi-user.target
