# products/fortress/systemd/fts-lte-collector.service
[Unit]
Description=Fortress LTE Usage Tracker (Watermark-based)
After=network.target
# Only run if the script exists and LTE interface is available
ConditionPathExists=/opt/hookprobe/fortress/devices/common/lte-usage-tracker.sh

[Service]
Type=oneshot
ExecStart=/opt/hookprobe/fortress/devices/common/lte-usage-tracker.sh collect
User=root
Group=root

# Minimal resources
Nice=19
IOSchedulingClass=idle

# Don't fail the boot if this fails
SuccessExitStatus=0 1

[Install]
WantedBy=multi-user.target
