# products/fortress/systemd/fts-device-status.service
[Unit]
Description=HookProbe Fortress Device Status Updater
Documentation=https://github.com/hookprobe/hookprobe
After=network.target openvswitch-switch.service fortress.service
# Only run if the script exists
ConditionPathExists=/opt/hookprobe/fortress/devices/common/device-status-updater.sh

[Service]
Type=oneshot
ExecStart=/opt/hookprobe/fortress/devices/common/device-status-updater.sh
User=root
Group=root

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fts-device-status

# Don't fail the boot if this fails
SuccessExitStatus=0 1

[Install]
WantedBy=multi-user.target
