# products/guardian/config/systemd/guardian-ssid-health.service
[Unit]
Description=Guardian SSID Health Check
Documentation=https://github.com/hookprobe/hookprobe
After=hostapd.service dnsmasq.service
Wants=hostapd.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/guardian-ssid-health.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=guardian-ssid-health

# Allow time for fixes
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
