# products/guardian/config/systemd/guardian-wifi-health.service
[Unit]
Description=Guardian WiFi WAN Health Check
Documentation=https://hookprobe.com/docs/guardian
After=network.target wpa_supplicant.service
Wants=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/guardian-wifi-health.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
